# delete_file.list
/usr/bin/droid/droid-hal-shutdown.sh
/usr/bin/droid/droid-hal-startup.sh
/usr/bin/droid/droid-init-done.sh
/usr/lib/systemd/system/actdead.target.wants/droid-hal-init.service
/usr/lib/systemd/system/droid-hal-init.service
/usr/lib/systemd/system/droid-late-start.service
/usr/lib/systemd/system/graphical.target.wants/droid-hal-init.service
/usr/lib/systemd/system/graphical.target.wants/droid-late-start.service
/etc/media_codecs.xml
/etc/media_codecs_ffmpeg.xml
/etc/media_codecs_google_audio.xml
/etc/media_codecs_google_telephony.xml
/etc/media_codecs_google_video.xml
/etc/media_codecs_performance.xml
/etc/media_profiles.xml
/etc/mixer_paths.xml
/etc/xtra_root_cert.pem
/etc/xtwifi.conf
/etc/sap.conf
/etc/sec_config
/etc/izat.conf
/etc/lowi.conf
/etc/acdbdata
/etc/gps.conf
/etc/audio_platform_info.xml
/usr/lib/systemd/user/lipstick.service.d/50-compositor.conf
/usr/lib/systemd/user/jolla-startupwizard-pre-user-session.service.d/50-compositor.conf
