# sparse/etc/systemd/user/lipstick.service.d/ofono-wait.conf
[Service]
ExecStartPre=/usr/bin/droid/lipstick-ofono-wait.sh
