# sparse/etc/systemd/system/ofono.service.d/ofono.conf
[Service]
ExecStartPost=/usr/bin/droid/ofono-delay.sh
TimeoutStartSec=60
