# sparse/etc/systemd/system/mce.service.d/backlight.conf
[Unit]
After=sys-devices-platform-backlight-backlight-backlight.device
Requires=sys-devices-platform-backlight-backlight-backlight.device
