# sparse/lib/systemd/system/enable-modem.service
[Unit]
Description=Enable the EG25 modem at boot

[Service]
Type=oneshot
ExecStart=/usr/bin/droid/enable-modem.sh

[Install]
WantedBy=basic.target
