# sparse/usr/lib/systemd/user/user-session.target.wants/unmute-sound-card.service
[Unit]
Description=Unmutes the A64 sound card at boot for PulseAudio

[Service]
Type=oneshot
ExecStart=/usr/bin/droid/unmute-sound-card.sh

[Install]
WantedBy=sound.target
